# pos-system/rootfs/etc/services.d/pos/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop POS System
# ==============================================================================

bashio::log.info "Stopping POS System..."
