# pos-system/rootfs/etc/services.d/pos/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# Read add-on configuration (prefer bashio, otherwise try /data/options.json, then env defaults)

get_json_value() {
  # args: key default
  local key="$1" default="$2"
  if [ -f /data/options.json ]; then
    local val
    val=$(sed -n 's/.*"'"${key}"'"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p' /data/options.json || true)
    if [ -n "${val}" ]; then
      printf '%s' "${val}"
      return 0
    fi
    # try numeric
    val=$(sed -n 's/.*"'"${key}"'"[[:space:]]*:[[:space:]]*\([0-9][0-9]*\).*/\1/p' /data/options.json || true)
    if [ -n "${val}" ]; then
      printf '%s' "${val}"
      return 0
    fi
  fi
  printf '%s' "${default}"
}

# Prefer /data/options.json (present in add-on runtime). Fall back to env defaults.
DB_HOST=$(get_json_value db_host core-mariadb)
DB_PORT=$(get_json_value db_port 3306)
DB_NAME=$(get_json_value db_name pos_db)
DB_USER=$(get_json_value db_user pos_user)
DB_PASSWORD=$(get_json_value db_password password)
ADMIN_USER=$(get_json_value admin_user admin)
ADMIN_PASSWORD=$(get_json_value admin_password admin123)
INGRESS_PATH=""

export DB_HOST DB_PORT DB_NAME DB_USER DB_PASSWORD ADMIN_USER ADMIN_PASSWORD INGRESS_PATH

echo "[POS] Starting POS System..."
echo "[POS] Database: ${DB_HOST}:${DB_PORT}/${DB_NAME}"
if [ -n "${INGRESS_PATH}" ]; then
  echo "[POS] Ingress path: ${INGRESS_PATH}"
fi

cd /app || exit 1
exec node server.js
